# instagram_summary.crontab
# m h dom mon dow command
0 15 * * * cd $HOME/instagram_summary && ./upload.sh